# scripts/systemd/sentinel-api.service
[Unit]
Description=Sentinel Core API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=__SENTINEL_HOME__
EnvironmentFile=/etc/sentinel/sentinel.env
# Loopback only — never reachable from outside this host. The only
# intended access path is the WHM plugin's CGI reverse proxy (see
# whm-plugin/), which itself only runs on this same host. See
# docs/architecture/decisions/0002-whm-plugin-internal-only-backend.md.
ExecStart=__SENTINEL_HOME__/.venv/bin/uvicorn sentinel.bootstrap:app --host 127.0.0.1 --port 8443
Restart=on-failure
RestartSec=5

# Deliberately no PrivateTmp=, ProtectHome=, ProtectSystem=, or
# DynamicUser= here. Every one of them would silently break Sentinel:
# PrivateTmp gives this unit a private, empty /tmp while forensics
# monitoring watches the real one; ProtectHome/ProtectSystem=strict block
# the account-home read/write access quarantine needs; DynamicUser
# conflicts with running as a fixed root user (required — see
# docs/architecture/decisions/0001-v1-deployment-privilege-and-mode.md).
# Do not add systemd sandboxing directives to this unit without reading
# that ADR first.

[Install]
WantedBy=multi-user.target
